# supervisor.conf
[program:_APP_NAME_]
process_name=%(program_name)s
command=php _DEPLOY_PATH_/current/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/_APP_NAME__horizon.log