# .deployer.example/supervisor.conf
[program:{{application}}]
process_name=%(program_name)s
command=php {{deploy_path}}/current/artisan horizon
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/{{application}}__horizon.log